# ansible/roles/kenban/files/kenban-host-agent.service
[Unit]
Description=Kenban Host Agent
After=docker.service
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
WorkingDirectory=/home/pi/kenban
User=pi

ExecStart=/usr/bin/python3 /home/pi/kenban/host_agent.py
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
